# SystemdFiles/server-backup-v2.timer
[Unit]
Description=server-backup-v2

[Timer]
Unit=server-backup-v2.service
OnCalendar=*-*-* 08:00:00

[Install]
WantedBy=timers.target